# roles/systemd_networkd/templates/reload-udev-offline.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Trigger reloading of udev while offline

[Timer]
OnActiveSec=5
Unit=reload-udev-offline.service
