# roles/systemd_networkd/templates/reload-udev-offline.service.j2
# {{ ansible_managed }}

[Unit]
Description=Coldplug All udev Devices (Offline)

[Service]
ExecStart={{ systemd_networkd_reload_udev_offline_script.dest }}
