# playbooks/templates/down-vpn-interface.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Trigger downing of the VPN network interface device

[Timer]
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=multi-user.target
