# playbooks/templates/up-vpn-interface.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Trigger upping of the VPN network interface device

[Timer]
OnCalendar=*-*-* 12:00:00

[Install]
WantedBy=multi-user.target
